# config/php/supervisor/test.conf
[program:test]
command=/usr/local/bin/php /code/bin/console myTest:supervisor
autorestart=true
stdout_logfile=/var/log/supervisor/test.log